# etc/network/if-up.d/iptables
#!/bin/bash
/sbin/iptables-restore < /etc/iptables.ipv4.nat
